# RNpicker/etc/ext/crontab/echo_time.sh
#!/bin/sh

# for testing a crontime job

date=`date`

echo "Hello world at $date"
